# deploy/jetson/go2-brain.service.in
[Unit]
Description=go2_local_brain AI autonomy GUI (auto-start)
After=network-online.target go2-ollama.service
Wants=network-online.target
Requires=go2-ollama.service

[Service]
Type=simple
User=@USER@
WorkingDirectory=@REPO@
Environment=HOME=@HOME@
EnvironmentFile=-@REPO@/.env
EnvironmentFile=-@HOME@/.go2/env.local
# Bind to all interfaces with a STABLE auth token from env so the laptop
# URL doesn't change every restart. Set GO2_GUI_TOKEN in ~/.go2/env.local
# (the installer generates one). If unset, --auth-token "" with empty value
# would be rejected, so the installer guarantees it is set.
ExecStart=@VENV@/bin/python -m go2_local_brain.ai_autonomy_gui --host 0.0.0.0 --port 8775 --auth-token ${GO2_GUI_TOKEN}
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
