# deploy/systemd/go2-local-brain.service
[Unit]
Description=Go2 local brain Jetson cockpit
After=network-online.target ollama.service
Wants=network-online.target ollama.service

[Service]
Type=simple
User=jetson
Group=jetson
WorkingDirectory=/home/jetson/robotics/go2_local_brain
EnvironmentFile=-/home/jetson/robotics/go2_local_brain/.env
ExecStart=/home/jetson/robotics/go2_local_brain/scripts/run_jetson_cockpit.sh
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
